# pptp/supervisord.conf
[supervisord]
nodaemon=true

[program:pptpd]
command=/run.sh
autorestart=true

