# pptp/supervisord.conf
[supervisord]
nodaemon=true

[program:pptpd]
command=/usr/sbin/pptpd -f
autorestart=true

