# mysql/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/bin/bash -c "mysqld --user=mysql --init-file=/init-mysql.sql"
autorestart=true

