# php-server/supervisord.conf
[supervisord]
nodaemon=true

[program:php5-fpm]
command=/usr/sbin/php5-fpm
autorestart=true

[program:nginx]
command=/server/tengine/sbin/nginx
autorestart=true


