# ssdb/supervisord.conf
[supervisord]
nodaemon=true

[program:ssdb]
command=/usr/local/ssdb/ssdb-server -s restart /usr/local/ssdb/ssdb.conf
autorestart=true


